# dwm/etc/xdg/autostart/slstatus.desktop
[Desktop Entry]
Type=Application
Exec=slstatus
Hidden=false
Name=slstatus for dwm
