# dwm/etc/xdg/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
Exec=udiskie -t
Hidden=false
Name=udiskie
