# base/etc/xdg/autostart/pywall.desktop
[Desktop Entry]
Type=Application
Exec=wal -R
Name=Wall Restore
